# issabel-contacts-uploader.service
[Unit]
Description=Issabel 4 Contacts Uploader Dashboard
After=network.target

[Service]
Type=simple
User=asterisk
WorkingDirectory=/opt/issabel-contacts-uploader
ExecStart=/usr/bin/python3 app.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
